# roles/nginx_exporter/templates/etc/systemd/system/prometheus-nginx-exporter.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus-nginx-exporter \
  -nginx.scrape-uri "$NGINX_SCRAPE_URI" \
  -web.listen-address "$LISTEN_ADDRESS" \
  $EXTRA_ARGS
