# roles/postgres_exporter/templates/etc/systemd/system/prometheus-postgres-exporter.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/prometheus-postgres-exporter.private
ExecStart=
ExecStart=/usr/bin/prometheus-postgres-exporter \
  --web.listen-address "$LISTEN_ADDRESS" \
  $EXTRA_ARGS
