# roles/fluent_bit/templates/etc/systemd/system/fluent-bit.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/fluent-bit/bin/fluent-bit \
  --config /etc/fluent-bit/fluent-bit.conf \
  --log_file /var/log/fluent-bit.log
