# roles/redis_exporter/templates/etc/systemd/system/prometheus-redis-exporter.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/prometheus-redis-exporter.private
ExecStart=
ExecStart=/usr/bin/prometheus-redis-exporter \
  -web.listen-address "$LISTEN_ADDRESS" \
  $EXTRA_ARGS
