# roles/grafana/templates/etc/systemd/system/grafana-server.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/grafana-server.private
