# roles/logrotate/templates/etc/systemd/system/logrotate.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=hourly
AccuracySec=5min
