# roles/node_exporter/templates/etc/systemd/system/prometheus-node-exporter.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus-node-exporter \
  --web.listen-address "$LISTEN_ADDRESS" \
  --collector.textfile.directory "$COLLECTOR_DIR" \
  $EXTRA_ARGS
