# roles/firefish/templates/etc/systemd/system/firefish.service.d/override.conf
[Unit]
After=postgresql.service
After=redis-server.service
