# roles/prometheus/templates/etc/systemd/system/prometheus.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus \
  --web.listen-address "$LISTEN_ADDRESS" \
  $EXTRA_ARGS
StandardOutput=file:/var/log/prometheus/prometheus.log
